# playbooks/apps/uptime/files/systemd/uptime-kuma.container
[Service]
Restart=always

[Container]
ContainerName=uptime-kuma
Image=docker.io/louislam/uptime-kuma:latest
Label=io.containers.autoupdate=registry
Volume=/${HOME}/data:/app/data:Z
PodmanArgs=--dns=9.9.9.9
Network=host

[Install]
WantedBy=default.target
